# docker/php/supervisord.conf
[supervisord]
logfile = /code/var/log/supervisord.stdout.log
user=root

[program:php-fpm]
command=php-fpm
user=root
autostart=true
autorestart=true
stderr_logfile=/code/var/log/php.stderr.log

[program:messenger-consume]
command=php /code/bin/console messenger:consume async --time-limit=3600
user=root
numprocs=2
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/code/var/log/queue.stdout.log
stderr_logfile=/code/var/log/queue.stderr.log
